# docker/medium/conf/php/supervisord.ini
[supervisord]
nodaemon=true

# [supervisorctl]
# serverurl=http://127.0.0.1:9001

# [inet_http_server]
# port=127.0.0.1:9001

[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autorestart=unexpected
exitcodes=0
stdout_logfile=/var/log/nginx.log
stderr_logfile=/var/log/nginx-error.log

[program:backend]
;command=/usr/local/bin/php -S 0.0.0.0:6132 -t /var/api /var/api/processor.php
command=/usr/local/bin/php -S 0.0.0.0:9000 /var/api/processor.php
autorestart=true
stdout_logfile=/var/log/backend.log
stderr_logfile=/var/log/backend-error.log
;stdout_logfile_maxbytes=0
;stderr_logfile_maxbytes=0
